# .config/systemd/user/gcalcli.timer
[Unit]
Description=gcalcli reminder timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=5m
Unit=gcalcli.service

[Install]
WantedBy=timers.target
